# cron/crontab
* * * * * /app/expire_tickets.sh >> /var/log/cron.log 2>&1

